# archive/systemd/user/fetch-pkg-list/fetch-pkg-list.service
[Unit]
Description=Regular Backup of Installed Packages
After=network.target

[Service]
ExecStart=%h/scripts/backup-pkg-list.sh
Type=oneshot
