# systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -w 40 -c 15 -d 5
